# backend/mvg-events.service
[Unit]
Description=Start MVG Event Bridge

[Service]
Type=simple
ExecStart=/srv/mvg-live/backend/mvg-live

[Install]
WantedBy=multi-user.target